# config/supervisor/gunicorn.ini
[program:sorting_visualizer_gunicorn]
command = /var/www/sorting-visualizer/config/gunicorn/start           ; Command to start app
user = spacer                                                          ; User to run as
stdout_logfile = /var/www/sorting-visualizer/logs/gunicorn_supervisor.log   ; Where to write log messages
redirect_stderr = true                                                ; Save stderr in the same log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8                       ; Set UTF-8 as default encoding
